# scripts/launchd/com.tanishk.driftedge.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tanishk.driftedge.poll</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tanishkyadav/dev/DriftEdge/scripts/driftedge-poll.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tanishkyadav/dev/DriftEdge</string>

    <!-- Start immediately on load and on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive: relaunch if it exits for any reason.
         ExitTimeOut=60 gives the process time to flush logs before SIGKILL. -->
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ExitTimeOut</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/tanishkyadav/dev/DriftEdge/logs/launchd-poll.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/tanishkyadav/dev/DriftEdge/logs/launchd-poll.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TOP_N</key>
        <string>20</string>
        <key>BOOK_INTERVAL_S</key>
        <string>30</string>
        <key>MARKET_REFRESH_S</key>
        <string>300</string>
    </dict>
</dict>
</plist>
